# IIoTGateway/iiotgateway.service
[Unit]
Description=IIoTGateway Service

[Service]
ExecStart=/opt/iiotgateway/rungateway.sh
Restart=always

[Install]
WantedBy=multi-user.target