# distro/board/iiotgateway/raspberrypi3-64/rootfs_overlay/usr/lib/systemd/system/gateway.service
[Unit]
Description=Gateway Service

[Service]
ExecStart=/usr/bin/gateway-start
Restart=always

[Install]
WantedBy=multi-user.target