# IIoTGateway/service/iiotgateway.service
[Unit]
Description=IIoTGateway Service

[Service]
ExecStart=/opt/iiotgateway/rungateway
Restart=always

[Install]
WantedBy=multi-user.target